# packaging/systemd/loncher.service
[Unit]
Description=loncher desktop control plane
Documentation=https://github.com/Mar2ianen/loncher

[Service]
Type=simple
ExecStart=%h/.local/bin/loncher daemon
EnvironmentFile=-%h/.config/loncher/environment
Restart=on-failure
RestartSec=2s
TimeoutStopSec=10s

[Install]
WantedBy=default.target
